# roles/service-wiki/templates/sync_wiki.timer.j2
[Unit]
Description=Sync Wiki data with Gogs every minute

[Timer]
OnBootSec=5min
OnUnitInactiveSec=1min

[Install]
WantedBy=timers.target