# roles/ff-supernode/templates/nodeinfo.cron.j2
* * * * * root {{supernode_freifunk_tools_dir}}/nodeinfo/nodeinfo.py | gzip | alfred -s 158
