# roles/ffmap-backend/templates/ffmap-backend.cron.j2
* * * * * root {{ffmap_backend_install_dir}}/backend.py -d {{ffmap_backend_data_dir}} -p {{ffmap_backend_prune_days}}
