# roles/internal/ffdo.mapserver/templates/alfred-master.service.j2
[Unit]
Description=Almighty Lightweight Fact Remote Exchange Daemon - Master

[Service]
ExecStart=/usr/local/sbin/alfred -m -i alfred0
Restart=on-failure

[Install]
WantedBy=multi-user.target

